# etc/cron.hourly/etckeeper-commit
#!/bin/sh


export LANG=C

cd /etc
[ -x `which etckeeper` ] && etckeeper commit "Hourly commit" >/dev/null 2>&1
exit 0
