# etc/cron.hourly/silent-auto-apt-get-upgrade
#!/bin/sh

export LANG=C

apt-get -qq -y -u update && apt-get -qq -y -u upgrade

